# mkl/root/lib/systemd/system/enable_mkl.service
[Unit]
Description=Custom Start and Shutdown
[Service]
Type=oneshot
ExecStart=/etc/mkl/mkl
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
